# systemd/panarchybotdev-watcher.path
[Path]
PathChanged=/usr/local/bin/panarchybotdev
[Install]
WantedBy=multi-user.target
