# systemd/panarchybot-watcher.path
[Path]
PathChanged=/usr/local/bin/panarchybot
[Install]
WantedBy=multi-user.target
